# dist/systemd/it100mqtt.service
[Unit]
Description=IT100 MQTT Bridge
After=ser2net.service emqx.service network.target
Requires=ser2net.service emqx.service
StartLimitIntervalSec=20

[Service]
Type=notify
User=root

ExecStart=/usr/local/bin/it100mqtt /usr/local/etc/it100mqtt.conf

TimeoutStartSec=10
NotifyAccess=all
Restart=always

[Install]
WantedBy=multi-user.target
